# crontab
*/30 * * * * cd /app && php process.php
